# ignition.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
 "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <!-- Unique identifier for the job -->
  <key>Label</key>
  <string>com.user_name.ignition</string>

  <!-- Run everything through bash so we can chain multiple commands -->
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>
      /usr/bin/python3 '/Users/user_name/path/to/your_script.py' ; \
      echo "Hello from your plist shell command!" ; \
      open -a 'Wireshark'
    </string>
  </array>

  <!-- Start automatically when you log in -->
  <key>RunAtLoad</key>
  <true/>

  <!-- Keep it alive (if process exits, launchd will restart it) -->
  <key>KeepAlive</key>
  <true/>

  <!-- Make sure it only runs in GUI login (needed for 'open -a') -->
  <key>LimitLoadToSessionType</key>
  <string>Aqua</string>
</dict>
</plist>
